# auto-scripts/files/auto-reboot-script
#!/bin/sh /etc/rc.common

START=99

start() {
    [ -f "/usr/bin/openwrt" ] && openwrt -e
}

stop() {
    # 停止操作（如果需要）
    echo "Service stopped."
}
