# auto-scripts/files/networkdetection
#!/bin/sh /etc/rc.common

START=99

start() {
    echo "Starting network check service..."
    /etc/init.d/check_network enable
}

stop() {
    echo "Stopping network check service..."
    /etc/init.d/check_network disable
}

check_network() {
    # 定义要检测的 DNS 服务器列表
    DNS_SERVERS=("8.8.8.8" "1.1.1.1" "119.29.29.29" "223.5.5.5")

    # 初始化网络状态标志
    network_up=0

    # 遍历 DNS 服务器列表
    for dns in "${DNS_SERVERS[@]}"; do
        # 尝试向每个 DNS 服务器发送 ping 请求
        ping -c 3 "$dns" > /dev/null 2>&1

        if [ $? -eq 0 ]; then
            echo "$(date): Network is up via $dns." >> /var/log/network_check.log
            network_up=1
            break  # 如果有一个 DNS 服务器响应，则退出循环
        else
            echo "$(date): Network check failed via $dns." >> /var/log/network_check.log
        fi
    done

    # 如果所有 DNS 服务器都无法响应，则认为网络不可用
    if [ $network_up -eq 0 ]; then
        echo "$(date): Network is down. Restarting network services..." >> /var/log/network_check.log
        /etc/init.d/network restart
        /etc/init.d/dnsmasq restart
    fi
}
