# auto-scripts/files/Restart_scripts
#!/bin/sh /etc/rc.common

START=99

start() {
    /etc/scripts_reboot
    openwrt -e
}

stop() {
    # 停止操作（如果需要）
    echo "Service stopped."
}
